# src/system-tests/coredump/systemtest-coredump.service
[Unit]
Description=Test a coredump

[Service]
ExecStart=/usr/bin/gnome-continuous-systemtest-coredump
